# hotspot_autostart.service
[Unit]
Description=Hotspot Autostart
After=network.target
[Service]
ExecStart=/usr/sbin/hotspot_startup.sh
[Install]
WantedBy=multi-user.target